# trunk/wasabi/sys/freevo.initscript
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting Freevo session"
	start-stop-daemon --start --quiet --background --make-pidfile \
		--pidfile /var/run/freevo-x.pid --exec /bin/su \
		-- - freevo -c "PATH=$PATH:/usr/X11R6/bin startx"
	eend $?
}

stop() {
	ebegin "Stopping Freevo session"
	start-stop-daemon --stop --quiet --pidfile /var/run/freevo-x.pid
	eend $?
}
